# systemd/briefing-daily.timer
[Unit]
Description=Run daily general morning briefing at 05:30

[Timer]
Unit=briefing@daily.service
OnCalendar=*-*-* 05:30:00
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
