# systemd/briefing-tech.timer
[Unit]
Description=Run tech morning briefing on Mondays at 06:00

[Timer]
Unit=briefing@tech.service
OnCalendar=Mon *-*-* 06:00:00
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
