# systemd/briefing-weekend.timer
[Unit]
Description=Run AI morning briefing on Wednesdays at 06:00

[Timer]
Unit=briefing@ai.service
OnCalendar=Wed *-*-* 06:00:00
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
